# legacy/Linux/configs/crontab
0 * * * * export DISPLAY=:0 && notify-send 'Hour pass!'

